# configs/system/proton-drive.service
[Unit]
Description=Mount Proton Drive using rclone
After=network-online.target

[Service]
Type=notify
ExecStart=/usr/bin/rclone mount proton: /home/lord/ProtonDrive --vfs-cache-mode=full --vfs-cache-max-age=1w 
ExecStop=/usr/bin/fusermount -uz /home/lord/ProtonDrive
Restart=on-failure
User=lord
Group=lord

[Install]
WantedBy=default.target

